# sources/meta-flexaudio/recipes-multimedia/flexaudio/files/flexaudio.service

[Unit]
Description=FlexAudio

[Service]
ExecStart=/usr/bin/flexaudio

[Install]
WantedBy=multi-user.target